# supervisor.conf
[program:uwsgi]
command = /home/docker/bin/uwsgi --ini /home/docker/uwsgi.ini

[program:nginx]
command = /usr/sbin/nginx
